# server-api/server-api.conf
[program:server-api]
command = python2.7 /home/smv/odoo-selhoz/server-api/server.py
autorestart = true
stderr_logfile = /var/log/server-api.err.log
stdout_logfile = /dev/null
